# config/auto_scheduler_task.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>Auto-Scheduler Báo Cáo Giao Hàng Telegram — Chạy ngầm gửi KRC/KSL/ĐÔNG/MÁT theo lịch cố định. Tự khởi động khi logon hoặc resume từ sleep.</Description>
  </RegistrationInfo>
  <Triggers>
    <!-- Trigger 0: Daily safety net — always start before 09:00 batch -->
    <CalendarTrigger>
      <StartBoundary>2026-05-30T08:45:00</StartBoundary>
      <Enabled>true</Enabled>
      <ScheduleByDay>
        <DaysInterval>1</DaysInterval>
      </ScheduleByDay>
    </CalendarTrigger>
    <!-- Trigger 1: Logon (reboot / first login) -->
    <LogonTrigger>
      <Enabled>true</Enabled>
      <Delay>PT3M</Delay>
    </LogonTrigger>
    <!-- Trigger 2: Resume from sleep / Modern Standby -->
    <SessionStateChangeTrigger>
      <Enabled>true</Enabled>
      <Delay>PT1M</Delay>
      <StateChange>ConsoleConnect</StateChange>
    </SessionStateChangeTrigger>
    <!-- Trigger 3: Unlock screen -->
    <SessionStateChangeTrigger>
      <Enabled>true</Enabled>
      <Delay>PT1M</Delay>
      <StateChange>SessionUnlock</StateChange>
    </SessionStateChangeTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <!-- No ExecutionTimeLimit — this is a long-running daemon -->
    <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>G:\My Drive\DOCS\transport_daily_report\tools\run-auto-scheduler.bat</Command>
      <WorkingDirectory>G:\My Drive\DOCS\transport_daily_report</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
